# config/systemd.launch.service
[Service]
Type=oneshot
Restart=no
RemainAfterExit=yes