# config/systemd.launch.service
[Service]
Type=oneshot
RemainAfterExit=yes